# deploy/systemd/valuechain-sec-acquisition.service
[Unit]
Description=ValueChain asynchronous SEC acquisition worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/ValueChain
EnvironmentFile=/home/pi/ValueChain/.env
ExecStart=/home/pi/ValueChain/.venv/bin/valuechain-acquire run-worker
Restart=on-failure
RestartSec=15s
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
TimeoutStopSec=2min

[Install]
WantedBy=default.target
