# deploy/systemd/valuechain-cvm-bulk-acquisition.service
[Unit]
Description=ValueChain Brazil CVM public bulk acquisition worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/ValueChain
EnvironmentFile=/home/pi/ValueChain/.env
ExecStart=/home/pi/ValueChain/.venv/bin/valuechain-global-acquire run-worker --source cvm_brazil
Restart=on-failure
RestartSec=30s
Nice=8
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStopSec=5min

[Install]
WantedBy=default.target
