# deploy/systemd/valuechain-cninfo-acquisition.service
[Unit]
Description=ValueChain low-rate CNINFO financial-report acquisition batch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/pi/ValueChain
EnvironmentFile=/home/pi/ValueChain/.env
ExecStart=/home/pi/ValueChain/.venv/bin/valuechain-global-acquire run-batch --source cninfo
Nice=10
IOSchedulingClass=idle
IOSchedulingPriority=7
TimeoutStartSec=30min
